# ubifs-root/miwifi_rd15_firmware_fd802_1.0.48.bin/squashfs-root/etc/hotplug.d/inited/03-port_service_arptables
#!/bin/sh

logger -p notice -t "hotplug.d" "03-port_service_arptables: run because of init done."

port_service reloadfw