# ubifs-root/miwifi_rd15_firmware_23a4f_1.0.68.bin/squashfs-root/etc/hotplug.d/ntp/10-miio-ntp.sh
#!/bin/sh

logger -t ntp "miio ntp hook"

if [ -n "$(pgrep miio_bind.sh)" ]; then
	logger -t miio_ntp "miio_bind.sh is running, skip it"
else
	miio_bind.sh -c
fi
