# ubifs-root/miwifi_rd15_firmware_176cc_1.0.65.bin/squashfs-root/usr/lib/opkg/info/anti-attack.list
/etc/init.d/anti_attack
/usr/sbin/anti_attack_adjust.sh
