# ubifs-root/miwifi_rd15_firmware_176cc_1.0.65.bin/squashfs-root/usr/lib/opkg/info/firewall_cpp.list
/usr/sbin/firewall_cpp
/usr/lib/libxtlock_intercept.so
/etc/init.d/firewall_cpp
/etc/config/firewall_cpp
