# ubifs-root/miwifi_rd15_firmware_176cc_1.0.65.bin/squashfs-root/usr/lib/opkg/info/kmod-qca-nss-ppe-ds.list
/etc/init.d/qca-nss-ppe-ds
/lib/modules/5.4.213/qca-nss-ppe-ds.ko
/etc/modules.d/51-qca-nss-ppe-ds
