# ubifs-root/miwifi_rd15_firmware_176cc_1.0.65.bin/squashfs-root/usr/lib/opkg/info/kmod-qca-nss-dp.list
/etc/init.d/qca-nss-dp
/usr/bin/edma_dump
/lib/modules/5.4.213/qca-nss-dp.ko
/etc/modules-boot.d/60-qca-nss-dp
/etc/modules.d/60-qca-nss-dp
/etc/config/qca_nss_dp
