# ubifs-root/miwifi_rd15_firmware_176cc_1.0.65.bin/squashfs-root/usr/lib/opkg/info/miparentalctl_v2.list
/etc/nginx/conf.d/8099.conf
/usr/sbin/mipctl_upg_ucfg
/etc/config/mipctl_user
/etc/init.d/mipctld
/usr/sbin/mipctl_upg.sh
/usr/sbin/mipctld
/lib/config_pre_ota/mipctl_dpi_ota.sh
/etc/config/mipctl_app
