# ubifs-root/miwifi_rd15_firmware_176cc_1.0.65.bin/squashfs-root/usr/lib/opkg/info/miwifi-logd.list
/etc/init.d/miwifi-logd
/usr/sbin/miwifi-logd
