# ubifs-root/miwifi_rd15_firmware_176cc_1.0.65.bin/squashfs-root/usr/lib/opkg/info/nfnlq-extension.list
/usr/sbin/nfnlq_userland
/etc/init.d/nfnlq_userland
/etc/config/nfnlq_userland
