# ubifs-root/miwifi_rd15_firmware_176cc_1.0.65.bin/squashfs-root/usr/lib/opkg/info/pluginmanager.list
/usr/share/pluginmanager/hosts.json
/etc/config/pluginmanager
/usr/sbin/pluginmanager
/etc/init.d/pluginmanager
/usr/sbin/pluginmanager_firewall
