# ubifs-root/miwifi_rd15_firmware_fd802_1.0.48.bin/squashfs-root/usr/lib/opkg/info/kmod-qca-wifi-lowmem-profile.list
/lib/modules/5.4.213/qca_ol.ko
/ini/internal/QCA9574_i.ini
/etc/init.d/load_cnss2
/lib/modules/5.4.213/ath_pktlog.ko
/ini/QCN9000.ini
/lib/wifi/qcacommands_ol_radio.xml
/ini/internal/QCA5018_i.ini
/ini/QCA5018.ini
/usr/sbin/thermaltool
/usr/sbin/exttool
/ini/internal/QCA6290_i.ini
/ini/internal/QCN9160_i.ini
/ini/internal/QCA5332_i.ini
/ini/QCN6122.ini
/usr/lib/libwlanfw_wmi.so
/ini/QCA8074.ini
/lib/wifi/launch_vap.sh
/ini/QCN9160.ini
/ini/QCA5332.ini
/etc/config/qcacfg80211
/etc/init.d/ssid_steering
/usr/sbin/cfg80211tool_mesh
/ini/internal/QCA8074V2_i.ini
/ini/internal/QCN9224_i.ini
/usr/sbin/apstats
/ini/QCA9984.ini
/ini/QCN9224.ini
/lib/modules/5.4.213/wifi_3_0.ko
/sawf/def_service_classes.json
/lib/wifi/tools_config
/lib/wifi/qcacommands_mesh.xml
/ini/QCA6018.ini
/ini/internal/global_i.ini
/lib/modules/5.4.213/mem_manager.ko
/ini/ftm.conf
/ini/QCA9574.ini
/etc/config/ssid-steering
/ini/internal/QCA6018_i.ini
/lib/modules/5.4.213/qdf.ko
/usr/sbin/assocdenialnotify
/etc/init.d/qcawifi-config-cmd
/usr/sbin/pktlogdump
/usr/sbin/80211stats
/usr/sbin/cfg80211tool
/lib/modules/5.4.213/monitor.ko
/lib/wifi/mount/umount_fw_partition.sh
/usr/lib/libqca_tools.so
/usr/lib/libwlanfw_halphy.so
/usr/sbin/radartool
/sawf/telemetry/sla.json
/ini/internal/QCA8074_i.ini
/lib/wifi/qcawifi_countrycode.txt
/lib/modules/5.4.213/umac.ko
/usr/sbin/wifistats
/usr/sbin/wifitool
/usr/sbin/wifi_hw_mode
/ini/internal/QCN9000_i.ini
/usr/sbin/acsdbgtool
/ini/internal/AR900B_i.ini
/usr/sbin/tx99tool
/usr/sbin/athstats
/usr/lib/libapstats.so
/lib/wifi/qcawificfg80211.sh
/usr/lib/libwlanfw_htt.so
/lib/wifi/mount/mount_fw_partition.sh
/lib/wifi/qca-wifi-modules
/etc/init.d/diag_socket_app
/ini/internal/QCN6432_i.ini
/sawf/service_classes.json
/lib/sbin/wifi_try
/ini/internal/QCN6122_i.ini
/sawf/telemetry/config.json
/usr/sbin/wlanconfig
/lib/wifi/qcacommands_vap.xml
/lib/wifi/debug/counters_wlan_prompt_hawkeye.sh
/lib/upgrade/keep.d/wifi_nss
/ini/global.ini
/usr/sbin/qca_gensock
/etc/hotplug.d/firmware/11-wifi-fw
/sawf/telemetry/sla_detect.json
/usr/lib/libwifistats.so
/lib/wifi/wifi-utils.sh
/lib/performance_configs.sh
/ini/QCN6432.ini
/usr/sbin/athstatsclr
/usr/sbin/ssidsteering
/usr/sbin/pktlogconf
/usr/sbin/wps_enhc
/ini/QCA6290.ini
/usr/sbin/wifi_try
/usr/sbin/athkey
