# ubifs-root/miwifi_rd15_firmware_fd802_1.0.48.bin/squashfs-root/usr/lib/opkg/info/kmod-qca-nss-ppe-bridge-mgr.list
/etc/modules.d/52-qca-nss-ppe-bridge-mgr
/etc/init.d/qca-nss-ppe-bridge-mgr
/lib/modules/5.4.213/qca-nss-ppe-bridge-mgr.ko
