# ubifs-root/miwifi_rd15_firmware_fd802_1.0.48.bin/squashfs-root/usr/lib/opkg/info/qti-license-pfm.list
/etc/init.d/license-pfm
/lib/upgrade/keep.d/license-pfm
/usr/sbin/license-pfm-upgrade.sh
