# ubifs-root/miwifi_rd15_firmware_fd802_1.0.48.bin/squashfs-root/usr/lib/opkg/info/sys_tester.list
/usr/sbin/memtester
/etc/init.d/systestctl
/usr/sbin/flashtestctl.sh
/usr/sbin/memtestctl.helper
