# ubifs-root/miwifi_rd15_firmware_fd802_1.0.48.bin/squashfs-root/etc/hotplug.d/iface/30-eth_iface_up.sh
#!/bin/sh
#logger -p notice -t "hotplug.d" "30-eth_iface_up.sh: run because of $INTERFACE $ACTION"

[ "$ACTION" = "ifup" ] && {
	NetMode=$(uci -q get xiaoqiang.common.NETMODE)
	if [ "$NetMode" == "whc_re" ]; then
		echo "eth_iface up,do topo changed check" >/dev/console
		ubus call topomon eth_up
	fi
}

[ "$ACTION" = "ifdown" ] && {
	NetMode=$(uci -q get xiaoqiang.common.NETMODE)
	if [ "$NetMode" == "whc_re" ]; then
		echo "eth_iface down,do topo changed check" >/dev/console
		ubus call topomon eth_down
	fi
}
