# ubifs-root/miwifi_rd15_firmware_23a4f_1.0.68.bin/squashfs-root/usr/lib/opkg/info/haohan_dpi.list
/etc/mipctl/_run.sh
/usr/lib/libsniper.so
/etc/init.d/haohan_dpi
/etc/mipctl/mipctl_app
/etc/mipctl/Sniper.klp
/usr/lib/iptables/libxt_miwifi-connhaohan.so
/usr/lib/libdpi_haohan.so
