# ubifs-root/miwifi_rd15_firmware_23a4f_1.0.68.bin/squashfs-root/usr/lib/opkg/info/local_gw_security.list
/etc/init.d/local_gw_security
/etc/config/local_gw_security
