# ubifs-root/miwifi_rd15_firmware_23a4f_1.0.68.bin/squashfs-root/usr/lib/opkg/info/lvm2.list
/sbin/lvmdiskscan
/etc/lvm/profile/lvmdbusd.profile
/sbin/vgcfgrestore
/sbin/pvs
/etc/lvm/profile/cache-smq.profile
/etc/lvm/profile/cache-mq.profile
/sbin/pvremove
/sbin/lvmsar
/etc/lvm/profile/command_profile_template.profile
/sbin/vgexport
/sbin/lvreduce
/sbin/lvrename
/sbin/pvmove
/sbin/dmsetup
/sbin/vgcfgbackup
/etc/init.d/lvm2
/sbin/vgcreate
/sbin/lvconvert
/sbin/lvremove
/sbin/vgmerge
/sbin/pvscan
/sbin/lvmconfig
/sbin/dmstats
/sbin/vgscan
/sbin/lvextend
/sbin/vgimport
/sbin/lvdisplay
/sbin/pvdisplay
/sbin/vgconvert
/sbin/pvresize
/sbin/pvchange
/sbin/vgck
/etc/lvm/lvm.conf
/sbin/vgmknodes
/sbin/vgreduce
/sbin/vgrename
/etc/lvm/profile/metadata_profile_template.profile
/sbin/lvresize
/sbin/lvchange
/etc/lvm/profile/thin-generic.profile
/sbin/vgremove
/sbin/vgdisplay
/sbin/lvscan
/sbin/pvcreate
/etc/lvm/profile/thin-performance.profile
/etc/lvm/lvmlocal.conf
/sbin/vgextend
/sbin/vgsplit
/sbin/lvm
/sbin/lvs
/sbin/lvcreate
/sbin/vgimportclone
/sbin/pvck
/sbin/lvmsadc
/sbin/vgs
/sbin/vgchange
