# ubifs-root/miwifi_rd15_firmware_23a4f_1.0.68.bin/squashfs-root/usr/lib/opkg/info/wififw_mount_script.list
/etc/init.d/wifi_fw_mount
/etc/hotplug.d/dump/00-q6dump
/etc/hotplug.d/dump_q6v5/00-q6dump
/etc/init.d/wifi_fw_done
